# EdenLinux/tmpl/packages/bootscripts/etc/init.d/modules
#!/bin/sh
### BEGIN INIT INFO
# Provides:          modules  
# Short-Description: module loading Startup Script
# Required-Start:	 
### END INIT INFO

. /etc/init.d/functions

# Assure that the kernel has module support.

case "$1" in
start)
	if [ -e /proc/ksyms -o -e /proc/modules ] ; then	
		echo -n "Loading modules: "
		
	    grep '^[^#]' /etc/modules |
		while read module args
		do
		    [ "$module" ] || continue
		    modprobe $module $args || :
			echo $module
		done
	
		check_status
	else
		echo "No module support in kernel!"
		exit 0
	fi
	;;
stop)
	;;
restart)
	start
	;;
*)
	echo "Usage: $0 {start|stop|restart}"
	exit 1
	;;
esac
