# EdenLinux/conf/packages/udev/udev
#!/bin/sh
### BEGIN INIT INFO
# Provides:          udev      
# Short-Description: Dynamic devices 
### END INIT INFO

. /etc/init.d/functions

case "$1" in
start)
	echo -n "Starting udev: "

	/bin/mount -n -t tmpfs -o mode=0755 udev /dev

	[ -d /lib/udev/devices ] && cp -af /lib/udev/devices/* /dev

	if [ ! -d /dev/shm ] ; then
		mkdir /dev/shm || die "failed to mkdir /dev/shm"
	fi

	chmod 1777 /dev/shm

	echo > /proc/sys/kernel/hotplug

	dmesg -n 2

	/sbin/udevd --daemon
	check_status
	
	mkdir -p /dev/.udev/queue

	echo -n "Waiting for devices to settle"
	/sbin/udevadm trigger
	/sbin/udevadm settle
	dmesg -n 7
	check_status
	;;
*)
	echo "Usage: $0 {start|stop|restart}"
	exit 1
esac