# EdenLinux/tmpl/packages/udev/udev
#!/bin/sh
### BEGIN INIT INFO
# Provides:          udev      
# Short-Description: Dynamic devices 
### END INIT INFO

. /etc/init.d/functions

case "$1" in
start)
	echo -n "Starting udev: "

	/bin/mount -n -t tmpfs -o mode=0755 udev /dev

	if [ -d /lib/udev/devices ] ; then
		cp -af /lib/udev/devices/* /dev
	fi
	
	if [ ! -d /dev/shm ] ; then
		mkdir /dev/shm || die "failed to mkdir /dev/shm"
	fi

	chmod 1777 /dev/shm

	echo > /proc/sys/kernel/hotplug

	dmesg -n 2
	
	if [ ! -d /run ] ; then
		mkdir /run || die "failed to mkdir /run"
	fi

	/sbin/udevd --daemon
	check_status
	
	mkdir -p /dev/.udev/queue

	echo -n "Waiting for devices to settle"
	/sbin/udevadm trigger
	/sbin/udevadm settle
	dmesg -n 7
	check_status
	
	echo -n "Mounting devpts: "
	mkdir /dev/pts
	mount -t devpts devpts /dev/pts 
	check_status
	;;
stop)
	echo -n "Stopping udev: "
	/sbin/udevadm control --exit
	check_status
	;;
*)
	echo "Usage: $0 {start|stop|restart}"
	exit 1
esac