# EdenLinux/conf/packages/bootscripts/etc/init.d/clock
#!/bin/sh
### BEGIN INIT INFO
# Provides:          syslog      
# Short-Description: Control the syslog daemon
### END INIT INFO

# Syslog Startup Script
#

. /etc/init.d/functions

SYSLOG_ROTATE_SIZE=65536

case "$1" in
start)
	echo -n "Setting system clock: "
	hwclock --hctosys --utc
	check_status
	;; 
stop)
	;;
restart)
	$0 stop
	$0 start
	;;
*)
	echo "Usage: $0 {start|stop|restart}"
	exit 1
esac