# EdenLinux/conf/packages/bootscripts/etc/init.d/mountroot
#!/bin/sh
### BEGIN INIT INFO
# Provides:          mountroot      
# Short-Description: Check and mount the root filesystem
# Required-Start:	 udev
### END INIT INFO

. /etc/init.d/functions
case "$1" in
start)
	echo "Mounting root file system in read-only mode."
	check_status
	mount -n -o remount,ro /
	
	echo "Checking root filesystem."
	/sbin/fsck -p -T /
	if [ "$?" -gt "1" ]; then
		echo "WARNING: Errors found while checking root filesystem."
		echo "You can login as root now, the system will reboot after logout."
		sulogin
		reboot
	fi
	
	echo -n "Remounting root in read-write mode: "
	mount -o remount,rw /
	check_status
	;;
stop)
	echo -n "Unmounting root"
	umount /
	check_status
	;;
*)
	echo "Usage: $0 {start|stop|restart}"
	exit 1
esac
