# EdenLinux/conf/packages/bootscripts/etc/init.d/syslog
#!/bin/sh
### BEGIN INIT INFO
# Provides:          syslog
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      
# Short-Description: Control the syslog daemon
### END INIT INFO

# Syslog Startup Script
#

. /etc/init.d/functions

SYSLOG_ROTATE_SIZE=65536

case "$1" in
start)
	echo -n "Starting syslogd: "
	syslogd -m 0 -s $SYSLOG_ROTATE_SIZE -L
	check_status
	echo -n "Starting klogd: "
	klogd
	check_status
	;;
stop)
	echo -n "Stopping klogd: "
	killall klogd
	check_status
	echo -n "Stopping syslogd: "
	killall syslogd
	check_status
	;;
restart)
	$0 stop
	$0 start
	;;
*)
	echo "Usage: $0 {start|stop|restart}"
	exit 1
esac